# apps/noodle/guadalinex-noodle/trunk/60-bluetooth.rules
# udev rules file for Bluetooth
SUBSYSTEM=="bluetooth", ACTION=="add", RUN+="/usr/bin/psuite"