# apps/noodle/guadalinex-noodle/trunk/60-ipaq.rules
# udev rules file for SynCE
PHYSDEVDRIVER=="ipaq", ACTION=="add", KERNEL=="ttyUSB*", RUN+="/usr/bin/psuite"
